# reboot/reboot.service
[Unit]
Description=Scheduled Daily Reboot

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --force reboot
